# hosts/hsb1/docker/restic-cron/hetzner/supervisor_restic.ini
[supervisord]
nodaemon=true

[program:cron]
startsecs=0
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
command=sh /usr/local/bin/start_cron.sh

[program:restic_backup]
startsecs=0
autostart=false
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/usr/local/bin/run_backup.sh

[program:restic_cleanup]
startsecs=0
autostart=false
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=sh /usr/local/bin/run_cleanup.sh

[program:restic_check]
startsecs=0
autostart=false
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=sh /usr/local/bin/run_check.sh
